# docs/crontabs/cron-job-list
# Load this with:
# sudo crontab -u ocpl cron-job-list

# m  h  dom mon dow command

# Main cronjobs controller
*/5  *  *   *   *    /home/ocpl/cron-defs/do-wget-url  util.sec/cron/run_cron.php  cron_main.html

# OKAPI cronjobs. Timeout of 1 second if okay, OKAPI keeps the process running
# in the background and has an error-reporting of its own. If something is really
# wrong (i.e. syntax error), then it will probably fail during the first second,
# making cron send a proper email to the admin. See OKAPI cronjobs dispatcher
# in okapi/cronjobs.php for more information.
*/5  *   *   *   *   wget -O - -q -t 1 https://opencaching.pl/okapi/cron5

# Database maintanance
10   6   *   *   *   /home/ocpl/cron-defs/do-mysql-flush

# Database backup
27   3   *   *   *   /home/ocpl/cron-defs/do-mysql-backup
