# docs/crontabs/do-mysql-flush
#!/bin/sh

/usr/bin/mysqladmin -u flusher -pflusherp flush-tables

