# doc/sql/crontabs/cron-job-list
# m  h  dom mon dow	command
*/30  *  *   *   *	/home/ocpl/cron-defs/do-wget-url   tpl/stdstyle/etc/write_newcaches.inc.php         write_newcaches.html
*/10 *  *   *   *	/home/ocpl/cron-defs/do-wget-url   tpl/stdstyle/etc/write_ratings.inc.php           write_ratings.html
*/15 *  *   *   *	/home/ocpl/cron-defs/do-wget-url   util.sec/email_sender/email_sender.php           email_sender.html
37   *  *   *   *	/home/ocpl/cron-defs/do-wget-url   util.sec/importcaches/import.php                 import.html
*/15 *  *   *   *	/home/ocpl/cron-defs/do-wget-url   util.sec/notification/run_notify.php             run_notify.html
0    *  *   *   *	/home/ocpl/cron-defs/do-wget-url   util.sec/publish_caches/run_publish.php          run_publish.html
0    23 *   *   *	/home/ocpl/cron-defs/do-wget-url   util.sec/search_index/fill_search_index.php      fill_search_index.html
*/20 *  *   *   *	/home/ocpl/cron-defs/do-wget-url   util.sec/watchlist/runwatch.php                  runwatch.html
30   4  *   *   *	/home/ocpl/cron-defs/do-wget-url   util.sec/auto_arch.php                           auto_arch.html
*/15 *  *   *   *	/home/ocpl/cron-defs/do-wget-url   util.sec/clearFakeVotes.php                      clearFakeVotes.html
*/5  *  *   *   *	/home/ocpl/cron-defs/do-wget-url   util.sec/geokrety/geokrety.class.php             geokrety.html
0    2  *   *   *	/home/ocpl/cron-defs/do-wget-url   util.sec/geokrety/geokrety.last.month.class.php  geokrety.last.month.html
18   *  *   *   *	/home/ocpl/cron-defs/do-wget-url   util.sec/repairScores.php                        repairScores.html
*/5  *  *   *   *       /home/ocpl/cron-defs/do-wget-url   util.sec/make_logo.php			    make_logo.html
# cache locations via NUTS Europe database
0    3  *   *   *	/home/ocpl/cron-defs/do-wget-url   util.sec/cache_locations/cache_location.class.php  cache_location.html
# Natura 2000 areas
30    3  *   *   *	/home/ocpl/cron-defs/do-wget-url   util.sec/cron/modules/cache_npa_areas.class.php  cache_npa_areas.html

#
#  MySQLa
0    *  *   *   *	php check_mysql_mem.php 	
10   6  *   *   *	/home/ocpl/cron-defs/do-mysql-flush
10   18 *   *   *	/home/ocpl/cron-defs/do-mysql-flush
#
# Backup database
17   5  *   *   *	/home/ocpl/cron-defs/do-mysql-backup
