# doc/sql/crontabs/do-copy-last-log.0
#!/bin/sh

LOGDIR=/home/ocpl/logs
LOGNAME=ocpl-access.log.1
NEWLOGNAME=ocpl-access.log.`date +%Y%m%d`

cp /var/log/apache2/$LOGNAME $LOGDIR
gzip $LOGDIR/$LOGNAME
mv $LOGDIR/$LOGNAME.gz $LOGDIR/$NEWLOGNAME.gz
chown ocpl:ocpl $LOGDIR/$NEWLOGNAME.gz
chmod 644 $LOGDIR/$NEWLOGNAME.gz

